# 10-ezxlayout.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="ezxlayout -a"
